# packaging/DEBIAN/postinst
#!/bin/bash
# Always exit on error
set -e

CONFIG_FILE="/etc/pifigo/config.yaml"
HOSTAPD_DEFAULT_FILE="/etc/default/hostapd"

# --- Helper function to read a value from the YAML config file ---
get_yaml_value() {
    # This simple grep/sed is good enough for this script.
    grep "^$1:" "$CONFIG_FILE" | sed "s/^$1: *//; s/\"//g" | tr -d '\r'
}

# --- Helper function to check for a working WLAN connection ---
check_wlan_connection() {
    WLAN_IFACE=$(get_yaml_value "wireless_interface")
    if [ -z "$WLAN_IFACE" ]; then
        echo "Could not read wireless_interface from config. Assuming offline."
        return 1
    fi

    echo "Checking for internet connectivity on interface: $WLAN_IFACE"
    if curl --interface "$WLAN_IFACE" --head --silent --fail --max-time 5 "http://www.google.com" > /dev/null; then
        echo "Success: The interface $WLAN_IFACE has a working internet connection."
        return 0
    else
        echo "Failure: The interface $WLAN_IFACE cannot reach the internet."
        return 1
    fi
}

# --- Main Installation Logic ---

echo "Running post-installation script for pifigo..."

# 1. Configure hostapd to use the config file we will generate.
# This follows the classic setup method.
echo "Configuring hostapd default settings..."
# Ensure the DAEMON_CONF line exists, then uncomment it and set the path.
if ! grep -q '^DAEMON_CONF=' "$HOSTAPD_DEFAULT_FILE"; then
    echo "DAEMON_CONF=\"/etc/hostapd/hostapd.conf\"" >> "$HOSTAPD_DEFAULT_FILE"
else
    # These two sed commands handle both commented and uncommented lines.
    sed -i 's|#DAEMON_CONF=.*|DAEMON_CONF="/etc/hostapd/hostapd.conf"|g' "$HOSTAPD_DEFAULT_FILE"
    sed -i 's|^DAEMON_CONF=.*|DAEMON_CONF="/etc/hostapd/hostapd.conf"|g' "$HOSTAPD_DEFAULT_FILE"
fi

# 2. Unmask dnsmasq, which is required for the hotspot.
if systemctl is-enabled --quiet dnsmasq.service | grep -q "masked"; then
    echo "Unmasking dnsmasq.service..."
    systemctl unmask dnsmasq.service
fi
echo "Ensuring dnsmasq service is enabled..."
systemctl enable dnsmasq.service

# 3. Reload systemd to make it aware of the new pifigo.service file.
echo "Reloading systemd daemon..."
systemctl daemon-reload

# 4. Check the system's network state.
if check_wlan_connection; then
    # SCENARIO A: INTERNET CONNECTION EXISTS
    echo "Existing Wi-Fi connection detected."
    systemctl enable pifigo.service
    echo "Installation complete. pifigo will start on the next boot if this connection fails."
else
    # SCENARIO B: NO INTERNET CONNECTION
    echo "No existing Wi-Fi connection detected. Starting pifigo in hotspot mode."
    # The pifigo service itself will now create the correct hotspot configs on its first run.
    systemctl enable --now pifigo.service
    echo "Installation complete. The pifigo hotspot should be active."
    echo "It is recommended to reboot the device to ensure a clean network state."
fi

exit 0
